# cron/vitographer.cron
* * * * * root /opt/vitographer/vitographer --log-dir /opt/vitographer/logs /opt/vitographer/conf.yml
