# supervisor/graphite.conf
[program:graphite_api]
command = uwsgi -y /etc/graphite-api.yaml
