# supervisor/vcontrold.conf
[program:vcontrol]
command = /opt/vcontrold/vcontrold -x vcontrold.xml -d /dev/ttyUSB0 -n -g -p 4000 --logfile /var/log/vcontrold.log
directory = /opt/vcontrold
startretries = 100
autorestart = true
